# shophosting-dir-backup.timer
Description=Daily backup timer for /opt/shophosting

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
