# shophosting-backup.service
[Unit]
Description=ShopHosting.io Daily Backup
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/shophosting/scripts/backup.sh
Environment=HOME=/root
Environment=XDG_CACHE_HOME=/root/.cache
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
